# deploy/trading-weekly-report.service
[Unit]
Description=AI Trading — weekly report
After=network.target

[Service]
Type=oneshot
User=trader
WorkingDirectory=/home/trader/ai-trading-system
EnvironmentFile=/home/trader/ai-trading-system/.env
ExecStart=/home/trader/ai-trading-system/.venv/bin/python -m src.run_weekly_report
