# deploy/trading-degen.timer
[Unit]
Description=Run degen sleeve every 15 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=15min
AccuracySec=30s

[Install]
WantedBy=timers.target
